# app/bin/crontab.cron.sample
50  2 * * * /home/analytics/imon/app/bin/cleartmp.sh

 * * * * *   /home/analytics/imon/app/bin/cron_m1.sh
20 * * * *  /home/analytics/imon/app/bin/cron_h1.sh
50 3 * * *  /home/analytics/imon/app/bin/cron_d1.sh
50 4 * * 1  /home/analytics/imon/app/bin/cron_w1.sh
50 5 * * *  /home/analytics/imon/app/bin/cron_mo1.sh

00 8 * * * /home/analytics/imon/app/bin/imon.sh newsmaker --granularity d1
05 8 * * 1 /home/analytics/imon/app/bin/imon.sh newsmaker --granularity w1
10 8 1 * * /home/analytics/imon/app/bin/imon.sh newsmaker --granularity m1
